# files/config/systemd/user/common.target
[Unit]
Description=Services that are to be started with or without presence of an X server
