# files/config/systemd/user/graphical.target
[Unit]
Description=Services to be started only in presence of an X server
# Requires=xorg.target
Conflicts=console.target
