# files/config/systemd/user/pyload.service
[Unit]
Description=PyLoad
Wants=network.target

[Service]
ExecStart=/usr/bin/pyLoadCore

[Install]
WantedBy=common.target
